# ansible/datomic/datomic_console.service.j2
# This file was automatically generated by an Ansible script
[Unit]
Name=datomic_console
Description=Datomic console web interface - Cognitect

[Service]
Type=simple
PIDFile={{ datomic_console_pid_file }}
Restart=always
RestartSec=30
ExecStartPre=-/bin/mkdir -p {{ datomic_pid_path }}
ExecStartPre=/bin/chown {{ datomic_user }} {{ datomic_pid_path }}
ExecStart={{ datomic_current_path }}/bin/console -p 8888 alias datomic:dev://{{ datomic_console_hostname }}:4334
User={{ datomic_user }}
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target