# ansible/datomic/datomic.service.j2
# This file was automatically generated by an Ansible script
[Unit]
Name=datomic
Description=Datomic database engine from Cognitect

[Service]
Type=simple
PIDFile={{ datomic_pid_file }}
Restart=always
RestartSec=30
ExecStartPre=-/bin/mkdir -p {{ datomic_pid_path }}
ExecStartPre=/bin/chown {{ datomic_user }} {{ datomic_pid_path }}
ExecStart={{ datomic_current_path }}/bin/transactor {{ datomic_config_file }}
User={{ datomic_user }}
PermissionsStartOnly=true

[Install]
WantedBy=multi-user.target